# templates/systemd/hermes-orchestrator.service
[Unit]
Description=Squad-Mind Hermes Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.hermes/bin/hermes-orchestrator
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
